# crontab
* * * * * root date >> /var/log/cron.log
# New line character required!